# deploy/rag-agent.service
[Unit]
Description=RAG Agent (FastAPI + BGE-M3 + DeepSeek)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{APP_USER}}
Group={{APP_USER}}
WorkingDirectory={{APP_DIR}}
EnvironmentFile={{APP_DIR}}/.env
ExecStart={{APP_DIR}}/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=30
LimitNOFILE=65536

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/data {{APP_DIR}}

[Install]
WantedBy=multi-user.target
